# packaging/protondisk-mount.service.in
[Unit]
Description=ProtonDisk — Proton Drive mounted as a disk
Documentation=https://github.com/mbeulens/ProtonDisk

[Service]
Type=simple
# @REPO@ is substituted with the checkout path by install.sh; %h is the user home.
# The daemon wrapper waits (event-driven, no polling) for the network, mounts,
# and re-mounts automatically when connectivity returns.
ExecStart=@REPO@/scripts/protondisk-mount-daemon %h/ProtonDisk
# SIGTERM (service stop / logout) makes libfuse tear the mount down; this is a
# best-effort backstop (the leading '-' ignores failure if already unmounted).
ExecStopPost=-/usr/bin/fusermount3 -uz %h/ProtonDisk
# Backstop only — the daemon self-supervises, so this rarely triggers.
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
